# install/providers/ansible/files/devclaw.service
[Unit]
Description=DevClaw - AI Agent for Tech Teams
After=network.target
Documentation=https://github.com/jholhewres/devclaw

[Service]
Type=simple
User=devclaw
Group=devclaw
WorkingDirectory=/home/devclaw
ExecStart=/usr/local/bin/devclaw serve
Restart=on-failure
RestartSec=5

EnvironmentFile=-/home/devclaw/.env

NoNewPrivileges=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true

LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
